# nanopore-sync.service
[Unit]
Description=Nanopore-Sync Service
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/env python -m nanopore_sync --source /path/to/input --destination /path/to/output
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target